# 01-cumbox.rules
ACTION=="add", ATTRS{idVendor}=="05ac", ATTRS{idProduct}=="1227", RUN+="/opt/cumbox/checkra1n"
